# os_mac/dwaglac.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>@LDNAME_LAC@</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>LoginWindow</string>
        <string>Aqua</string>
	</array>
	<key>ProgramArguments</key>
    <array>
    	<string>@PATH_DWA@/native/DWAgentSvc.app/Contents/MacOS/DWAgentSvc</string>
		<string>@PATH_DWA@</string>
		<string>lac</string>
		<string>@EXE_NAME@</string>
	</array>
	<key>ProcessType</key>
    <string>Interactive</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
    <true/>
    <key>ExitTimeOut</key>
    <integer>10</integer>
</dict>
</plist>
